# api_agent/debian-package/wireguardbroker_1.0/etc/systemd/system/wireguardbroker.service
[Unit]
Description=Wireguard API Broker
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/wireguardbroker

[Install]
WantedBy=multi-user.target